# crontab
* * * * * echo "Running cloudupload" && cloudupload
* * * * * echo "Running rmlocal" && rmlocal
